# lab09/mephi.service
[Unitt]
Decription="MEPhI hello service"
After=sshd.service
[Service]
ExecStart=/usr/bin/printf "Hello, MEPhI!"
[Install]
WantedBy=multi-user.target